# ros/supervisord.conf
[supervisord]
nodaemon=true

; abilita `supervisorctl` da dentro il container: usato dalla suite di test
; (test/conftest.py) per mettere in pausa la simulazione ROS/Gazebo durante
; i test (pesante in CPU, non serve ai test che usano il generatore
; sintetico) e per rimetterla su alla fine.
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:vncserver]
command=/usr/bin/tigervncserver -localhost no -geometry 1280x800 -depth 24 -SecurityTypes None --I-KNOW-THIS-IS-INSECURE :1 -fg
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:novnc]
command=websockify --web=/usr/share/novnc/ 6080 localhost:5901
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:roscore]
command=/bin/bash -lc "source /opt/ros/noetic/setup.bash && roscore"
autorestart=true
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Servizio di controllo del generatore sintetico (Passo 11 estensione): la
; dashboard lo chiama via HTTP, non serve lanciarlo a mano dopo il boot.
[program:generator_service]
command=/bin/bash -lc "cd /opt/shf/generator && python3 generator_service.py"
autorestart=true
priority=40
redirect_stderr=true
stdout_logfile=/tmp/generator_service.log
stdout_logfile_maxbytes=10MB

; Servizio di controllo della flotta reale (Passo 14): ponte HTTP -> topic ROS
; per iniezione guasti live e assegnazione missioni (kafka_bridge.py /
; graph_navigator.py). E' anche un nodo rospy (a differenza di
; generator_service.py), serve roscore gia' su -- da qui lo "sleep".
[program:fleet_control_service]
command=/bin/bash -lc "sleep 5 && source /root/catkin_ws/devel/setup.bash && cd /root/catkin_ws/src/shf_bringup/scripts && python3 -u fleet_control_service.py"
autorestart=true
priority=45
redirect_stderr=true
stdout_logfile=/tmp/fleet_control_service.log
stdout_logfile_maxbytes=10MB

; Servizio di controllo dei run di valutazione sperimentale (Passo 13,
; estensione "risultati live"): la dashboard avvia run_effectiveness.py /
; run_efficiency.py on-demand invece di leggere PNG pre-generati.
[program:eval_service]
command=/bin/bash -lc "cd /opt/shf/eval && python3 -u eval_service.py"
autorestart=true
priority=45
redirect_stderr=true
stdout_logfile=/tmp/eval_service.log
stdout_logfile_maxbytes=10MB

; Simulazione multi-robot di riferimento (Passo 5): parte in automatico cosi'
; la dashboard mostra subito una flotta reale con un solo `docker compose
; up`. roslaunch attende/riprova da solo la connessione al master ROS, non
; serve un controllo esplicito di readiness su roscore qui.
[program:sim_multi_robot]
command=/bin/bash -lc "sleep 5 && source /root/catkin_ws/devel/setup.bash && roslaunch shf_bringup sim_multi_robot.launch"
autorestart=true
priority=50
startsecs=15
redirect_stderr=true
stdout_logfile=/tmp/sim.log
stdout_logfile_maxbytes=10MB
