# ros/supervisord.conf
[supervisord]
nodaemon=true

; abilita `supervisorctl` da dentro il container: usato dalla suite di test
; (test/conftest.py) per mettere in pausa la simulazione ROS/Gazebo durante
; i test (pesante in CPU, non serve ai test che usano il generatore
; sintetico) e per rimetterla su alla fine.
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:vncserver]
command=/usr/bin/tigervncserver -localhost no -geometry 1280x800 -depth 24 -SecurityTypes None --I-KNOW-THIS-IS-INSECURE :1 -fg
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:novnc]
command=websockify --web=/usr/share/novnc/ 6080 localhost:5901
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:roscore]
command=/bin/bash -lc "source /opt/ros/noetic/setup.bash && roscore"
autorestart=true
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Servizio di controllo del generatore sintetico (Passo 11 estensione): la
; dashboard lo chiama via HTTP, non serve lanciarlo a mano dopo il boot.
[program:generator_service]
command=/bin/bash -lc "cd /opt/shf/generator && python3 generator_service.py"
autorestart=true
priority=40
redirect_stderr=true
stdout_logfile=/tmp/generator_service.log
stdout_logfile_maxbytes=10MB

; Servizio di controllo della flotta reale (Passo 14): ponte HTTP -> topic ROS
; per iniezione guasti live e assegnazione missioni (kafka_bridge.py /
; graph_navigator.py). E' anche un nodo rospy (a differenza di
; generator_service.py), serve roscore gia' su -- da qui lo "sleep".
[program:fleet_control_service]
command=/bin/bash -lc "sleep 5 && source /root/catkin_ws/devel/setup.bash && cd /root/catkin_ws/src/shf_bringup/scripts && python3 -u fleet_control_service.py"
autorestart=true
priority=45
redirect_stderr=true
stdout_logfile=/tmp/fleet_control_service.log
stdout_logfile_maxbytes=10MB

; Servizio di controllo dei run di valutazione sperimentale (Passo 13,
; estensione "risultati live"): la dashboard avvia run_effectiveness.py /
; run_efficiency.py on-demand invece di leggere PNG pre-generati.
[program:eval_service]
command=/bin/bash -lc "cd /opt/shf/eval && python3 -u eval_service.py"
autorestart=true
priority=45
redirect_stderr=true
stdout_logfile=/tmp/eval_service.log
stdout_logfile_maxbytes=10MB

; Simulazione multi-robot di riferimento (Passo 5). Dal 2026-07-21 NON parte
; piu' in automatico (autostart=false): l'utente vuole controllare dalla
; dashboard quando farla partire e con quale scala (small/large, Passo 14) --
; fleet_control_service.py la avvia/ferma con `supervisorctl start/stop
; sim_multi_robot` (stesso meccanismo gia' usato da test/conftest.py per
; mettere in pausa la simulazione durante i test, invariato). La scala e'
; passata via un file marcatore (/tmp/shf_scale, scritto da
; fleet_control_service.py prima di avviare) perche' il comando di un
; programma supervisord e' statico e non puo' ricevere argomenti dinamici ad
; ogni `start`; se il file manca (mai scelta una scala) il default e' small.
[program:sim_multi_robot]
command=/bin/bash -lc "sleep 5 && source /root/catkin_ws/devel/setup.bash && SCALE=$(cat /tmp/shf_scale 2>/dev/null || echo small) && roslaunch shf_bringup sim_multi_robot.launch scale:=$SCALE"
autostart=false
autorestart=true
priority=50
startsecs=15
redirect_stderr=true
stdout_logfile=/tmp/sim.log
stdout_logfile_maxbytes=10MB
