# ros/supervisord.conf
[supervisord]
nodaemon=true

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:vncserver]
command=/usr/bin/tigervncserver -localhost no -geometry 1280x800 -depth 24 -SecurityTypes None --I-KNOW-THIS-IS-INSECURE :1 -fg
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:novnc]
command=websockify --web=/usr/share/novnc/ 6080 localhost:5901
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:roscore]
command=/bin/bash -lc "source /opt/ros/noetic/setup.bash && roscore"
autorestart=true
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:generator_service]
command=/bin/bash -lc "cd /opt/shf/generator && python3 generator_service.py"
autorestart=true
priority=40
redirect_stderr=true
stdout_logfile=/tmp/generator_service.log
stdout_logfile_maxbytes=10MB

[program:fleet_control_service]
command=/bin/bash -lc "sleep 5 && source /root/catkin_ws/devel/setup.bash && cd /root/catkin_ws/src/shf_bringup/scripts && python3 -u fleet_control_service.py"
autorestart=true
priority=45
redirect_stderr=true
stdout_logfile=/tmp/fleet_control_service.log
stdout_logfile_maxbytes=10MB

[program:eval_service]
command=/bin/bash -lc "cd /opt/shf/eval && python3 -u eval_service.py"
autorestart=true
priority=45
redirect_stderr=true
stdout_logfile=/tmp/eval_service.log
stdout_logfile_maxbytes=10MB

[program:sim_multi_robot]
command=/bin/bash -lc "sleep 5 && source /root/catkin_ws/devel/setup.bash && SCALE=$(cat /tmp/shf_scale 2>/dev/null || echo small) && RUN_ID=$(cat /tmp/shf_run_id 2>/dev/null || echo '') && roslaunch shf_bringup sim_multi_robot.launch scale:=$SCALE run_id:=$RUN_ID"
autostart=false
autorestart=true
priority=50
startsecs=15
redirect_stderr=true
stdout_logfile=/tmp/sim.log
stdout_logfile_maxbytes=10MB
